# openrc/docker-compose-service-template
#!/sbin/openrc-run

name="$RC_SVCNAME"
description="OpenRC (template) script for managing a docker-compose project as a service"

command="docker compose --project-directory ${COMPOSE_PROJECT_DIRECTORY:-/srv/docker/$name}"

# extra_commands="help"
# help() {
# 	einfo "Usage: suppose a Gitlab application packaged for docker-compose..."
# 	einfo " - Drop the docker-compose.yaml file under /srv/docker/gitlab"
# 	einfo " - Create the configuration file: /etc/conf.d/gitlab"
# 	einfo " - Set inside the required variable as follow: COMPOSE_PROJECT_DIRECTORY=\"/srv/docker/gitlab\""
# 	einfo " - Create a symlink in /etc/conf.d that points to this file nammed gitlab"
# }

depend() {
	need docker
}

start_pre() {
	# if the service is not run using a symlink
	if [ ! -L "$1" ]; then
		eerror "This file is a template and is not intended to be run directly."
		eerror "See the 'help' section of this service for how to use it."
		exit 1
	fi

	[ -f "/etc/conf.d/$name" ] && . /etc/conf.d/"$name"
}

start() {
	ebegin "Starting docker compose: $name"
	$command $command_args up -d
	eend $?
}

stop() {
	ebegin "Stopping docker compose: $name"
	$command $command_args down
	eend $?
}

status() {
	ebegin "Informations about docker compose: $name"
	$command $command_args ps
	$command $command_args logs
}
